# config/captive-portal.service
[Unit]
Description=Captive Portal Network Server
Documentation=https://github.com/rayscosta/captive-portal-network
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=captive-portal
Group=captive-portal
WorkingDirectory=/opt/captive-portal
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=captive-portal

# Variáveis de ambiente
Environment=NODE_ENV=production
Environment=PORT=3000

# Limites de recursos
LimitNOFILE=65536

# Segurança
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/captive-portal/database.sqlite

[Install]
WantedBy=multi-user.target
